# systemd-root/system/hosts-update.service
[Unit]
Description=Update hosts file.

[Service]
Type=oneshot
ExecStart=/usr/bin/hosts-update

[Install]
WantedBy=default.target
