# systemd-root/system/pacman-sync.service
[Unit]
Description=Sync pacman repositories

[Service]
Type=oneshot
ExecStart=/bin/pacman -Sy

[Install]
WantedBy=default.target
