# systemd-root/etc/systemd/system/pacman-sync.timer
[Unit]
Description=Sync pacman repositories

[Timer]
OnBootSec=3m
OnUnitActiveSec=1h
Unit=pacman-sync.service

[Install]
WantedBy=default.target
